# supervisord.conf
[supervisord]
nodaemon=true 
loglevel = debug
redirect_stderr=true
redirect_stdout=true

[program:awslogsd]
command=/usr/sbin/awslogsd
