# supervisord.conf
[supervisord]
nodaemon=true 
[program:awslogsd]
command=/usr/sbin/awslogsd
redirect_stderr=true

