# supervisord.conf
[supervisord]
nodaemon=true 
loglevel = debug

[program:awslogsd]
command=/usr/sbin/awslogsd
